# opt/etc/cron.hourly/rrdstorm_graph.sh
#!/bin/sh

rrdstorm.sh graph 0 1 2 3 4 5
